# guacamole.conf
[program:tomcat7]
command=/usr/share/tomcat7/bin/catalina.sh run
environment=CATALINA_BASE="/var/lib/tomcat7"
user=tomcat7

[program:guacd]
command=/usr/local/sbin/guacd -f
user=root

